# internet-logger.service
[Unit]
Requires=network.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=
ExecStart=/home/pi/Code/internet-connectivity/internet.sh
